# roles/common/files/update-motd.d/10-hostname
#!/bin/bash

# colors
undim="\e[0m"
yellow="\e[1;93m"

printf "$yellow"
/usr/bin/env figlet "$(hostname)"
printf "$undim"
